# Addons/LilacTV/repo/script.updateShortCut/script.updateShortCut_bk/resources/Seebo.cron
#Check Internet connection for every minute
* * * * * /etc/seebo/CheckInternet.sh

#Check the softeware to update for every hours at 0 minute
0 * * * * kodi-send --action="RunScript(/storage/.kodi/addons/script.updateShortCut/UpdateSkin.py)"

#Reset Active Flag for every day at 5 AM
5 19 * * * python /storage/.kodi/addons/script.updateShortCut/dbHandleRun.py
