# Addons/LilacTV/repo/script.updateShortCut/script.updateShortCut/resources/Seebo.cron
#Check Internet connection for every minute
* * * * * /etc/seebo/CheckInternet.sh

#Check the softeware to update for every hours at 0 minute
0 * * * * kodi-send --action="RunScript(/storage/.kodi/addons/script.updateShortCut/UpdateSkin.py)"
